# not-stowable/udev/rules.d/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", HOTPLUG=="1", RUN+="/home/wolfe/.i3/hotplug_monitor.sh"
